# immich-backup/immich-backup.timer
[Unit]
Description=Weekly Immich backup timer

[Timer]
OnCalendar=Wed 14:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
